# docker/production/gunicorn.conf
[program:gunicorn]
command=/start_app
directory=/
user=django
autostart=true
autorestart=true
redirect_stderr=true
stdout_logfile = /logs/gunicorn.log